# systemd/system/nftwifi-init.service
[Unit]
Description=initialise nftwifi runtime directory

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/nftwifi-init
