# systemd/system/wifi-kill.service
[Unit]
Description=stop all wireless units on all interfaces

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=no
